# crontab
# Gajana schedule (times are container-local, TZ=Asia/Kolkata).
# Output goes to supercronic's stdout -> Docker logs -> Dozzle.

# Daily 06:00 IST: fetch new statements + update the log, then push a single
# Uptime-Kuma health heartbeat (see src/monitor.py).
0 6 * * * cd /app && python main.py --daily

# Weekly Sunday 07:00 IST: back up Google Sheets to the local SQLite DB.
0 7 * * 0 cd /app && python main.py --backup-db

# Monthly 5th 08:00 IST: split the previous month's Google payslip into the
# salary sheet + the ledger (no arg -> previous calendar month; idempotent).
0 8 5 * * cd /app && python run_salary_splitter.py

# Every minute: poll Telegram for cash entries and append them to the Cash
# Transactions sheet (no-op until the bot is configured).
* * * * * cd /app && python run_telegram_bot.py
